# elixir/deploy/stelgano.service
# Systemd unit for the sTELgano Phoenix release.
#
# Install once on the droplet:
#   sudo cp deploy/stelgano.service /etc/systemd/system/stelgano.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable stelgano
#
# The deploy workflow stops this service, swaps in the new release
# under /opt/stelgano/current, and starts it again.

[Unit]
Description=sTELgano (Phoenix release)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=deploy
Group=deploy
WorkingDirectory=/opt/stelgano/current

# Env vars (SECRET_KEY_BASE, DATABASE_URL, PHX_HOST, ADMIN_PASSWORD, …)
# live in /opt/stelgano/.env — mode 0600, owned by the deploy user.
EnvironmentFile=/opt/stelgano/.env

ExecStart=/opt/stelgano/current/bin/stelgano start
ExecStop=/opt/stelgano/current/bin/stelgano stop

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/stelgano

[Install]
WantedBy=multi-user.target
